# config/sidekiq-okurite.service
# Sidekiq systemd サービスファイル（okurite専用）
# 設置先: /etc/systemd/system/sidekiq-okurite.service
#
# 操作コマンド:
#   sudo systemctl start sidekiq-okurite    （起動）
#   sudo systemctl stop sidekiq-okurite     （停止）
#   sudo systemctl restart sidekiq-okurite  （再起動 ※git pull後はこれを実行）
#   sudo systemctl status sidekiq-okurite   （状態確認）
#   journalctl -u sidekiq-okurite -f        （ログ確認）

[Unit]
Description=sidekiq for okurite
After=syslog.target network.target redis.target

[Service]
Type=simple

WorkingDirectory=/home/smart/okurite
User=smart
Group=smart

Environment=MALLOC_ARENA_MAX=2

ExecStart=/bin/bash -lc 'bundle exec sidekiq -C config/sidekiq.yml'
ExecReload=/bin/kill -TSTP $MAINPID

RestartSec=1
Restart=on-failure

SyslogIdentifier=sidekiq-okurite
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
